# Data/deadman@.timer
[Unit]
Description=Dead men tell tales after 1 week

[Timer]
OnCalendar=%i

[Install]
WantedBy=timers.target
